# athinfod.service
[Unit]
Description=Athinfod Application

[Service]
User=root
Group=root
ExecStart=/usr/bin/athinfod
StandardInput=socket
